# bundle_specs/control_center_bundle_v1_6_6/systemd/pi-control-center.service
[Unit]
Description=Pi Control Center
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/mnt/nas/homelab/control-center/app
Environment=PYTHONUNBUFFERED=1
Environment=TAILSCALE_FQDN=pi-nas.taild4713b.ts.net
ExecStart=/mnt/nas/homelab/control-center/venv/bin/python -m uvicorn app:app --host 127.0.0.1 --port 9000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
