# bundle_specs/control_center_bundle_v1_7_0/payload/bundle_specs/control_center_bundle_v1_7_0/payload/bundle_specs/control_center_bundle_v1_7_0/payload/bundle_specs/control_center_bundle_v1_7_0/payload/bundle_specs/control_center_bundle_v1_7_0/payload/bundle_specs/control_center_bundle_v1_7_0/payload/recovery/self_recovery_v3_4/systemd/homelab-self-heal-v3_4.timer
[Unit]
Description=Run homelab self-healing sweep v3.4 every 2 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=2min
Unit=homelab-self-heal-v3_4.service

[Install]
WantedBy=timers.target
