# systemd/raspi-homelab-watchdog.service.example
[Unit]
Description=Raspi Homelab Watchdog
