# systemd/raspi-homelab-python-framework.service.example
[Unit]
Description=Raspi Homelab Python Framework Control Center
