# bundle_specs/control_center_bundle_v1_7_0/payload/bundle_specs/control_center_bundle_v1_7_0/payload/bundle_specs/control_center_bundle_v1_7_0/payload/bundle_specs/control_center_bundle_v1_7_0/payload/bundle_specs/control_center_bundle_v1_7_0/payload/bundle_specs/control_center_bundle_v1_7_0/payload/bundle_specs/control_center_bundle_v1_7_0/payload/bundle_specs/control_center_bundle_v1_7_0/payload/bundle_specs/control_center_bundle_v1_7_0/payload/bundle_specs/control_center_bundle_v1_7_0/payload/bundle_specs/control_center_bundle_v1_7_0/payload/bundle_specs/control_center_bundle_v1_7_0/payload/bundle_specs/control_center_bundle_v1_7_0/payload/bundle_specs/control_center_bundle_v1_7_0/payload/bundle_specs/control_center_bundle_v1_7_0/payload/bundle_specs/control_center_bundle_v1_7_0/payload/recovery/self_recovery_v3_4/systemd/homelab-self-heal-v3_4.service
[Unit]
Description=Homelab self-healing sweep v3.4

[Service]
Type=oneshot
ExecStart=/usr/local/bin/homelab-self-heal-v3_4
User=root
Environment=DOCKER_ROOT_DIR=/mnt/nas/homelab/docker
Environment=CLOUDFLARED_IMAGE=cloudflare/cloudflared:2026.1.2
Environment=ALLOW_REBOOT=1
