# plugins/files/docker/media-ingest.service
[Unit]
Description=NAS Media Ingest (Incoming to Library)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
ExecStart=/usr/local/bin/media_ingest_watch.sh
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
